# systemd/user/vodik.target
[Unit]
Description=Stuff for blackb1rd
After=xorg.target

[Install]
Alias=default.target
